# build/services/php_fpm.conf
[program:php-fpm]
command=php-fpm -F
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
