# build/dev/services/vite.conf
[program:vite]
command=npm run dev
directory=/var/app/www/frontend
priority=500
user=app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
