# build/prod/services/cron.conf
[program:cron]
command=supercronic /etc/cron.d/app
priority=900
user=app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
