# build/supervisord.conf.tmpl
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord/supervisord.pid

[program:php-fpm]
command=php-fpm -F
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

{{if eq .Env.APPLICATION_ENV "development"}}
[program:vite]
command=npm run dev
directory=/var/app/www/frontend
priority=500
user=app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
{{end}}

[program:caddy]
command=caddy run --config /etc/Caddyfile
priority=900
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:cron]
command=supercronic /var/app/crontab
priority=900
user=app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
